# deploy/heagent.service
# ============================================================
# HeAgent Systemd 服务单元
# 安装: sudo cp heagent.service /etc/systemd/system/
#        sudo systemctl daemon-reload
#        sudo systemctl enable --now heagent
# ============================================================

[Unit]
Description=HeAgent - Self-improving AI Agent
Documentation=https://github.com/your-org/heagent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=heagent
Group=heagent

# 工作目录
WorkingDirectory=/opt/heagent

# 环境变量文件
EnvironmentFile=/opt/heagent/.env.production

# 可执行命令（交互模式）
ExecStart=/opt/heagent/.venv/bin/python -m heagent

# 重启策略
Restart=always
RestartSec=10

# 资源限制
LimitNOFILE=65536
LimitNPROC=1024
CPUQuota=100%
MemoryMax=1G

# 安全加固
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
CapabilityBoundingSet=
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# 日志
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
